# btkbdd/btkbdd@.service
[Unit]
Requires=bluetooth.service
Documentation=man:btkbdd(1)

[Service]
ExecStart=-/usr/sbin/btkbdd /dev/input/%I -c /var/lib/btkbdd/%I.cable
