# evmuxd/evmuxd@.service
[Unit]
Documentation=man:evmuxd(1)

[Service]
ExecStart=-/usr/sbin/evmuxd /dev/input/%I
